# ansible/roles/sse-proxy/files/openclaw-sse-proxy.service
[Unit]
Description=OpenClaw SSE keepalive reverse proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debian
Group=debian
WorkingDirectory=/home/debian/.openclaw/sse-proxy
ExecStart=/usr/bin/env node /home/debian/.openclaw/sse-proxy/server.js
Restart=on-failure
RestartSec=2
StartLimitBurst=10
StartLimitIntervalSec=60

# Default env (override here or via systemctl edit if needed)
Environment=SSE_PROXY_LISTEN_HOST=127.0.0.1
Environment=SSE_PROXY_LISTEN_PORT=8444
Environment=SSE_PROXY_UPSTREAM_HOST=127.0.0.1
Environment=SSE_PROXY_UPSTREAM_PORT=8443
Environment=SSE_PROXY_KEEPALIVE_MS=15000

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/debian/.openclaw/sse-proxy
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-sse-proxy

[Install]
WantedBy=multi-user.target
