# ansible/roles/gateway-proxy/files/openclaw-gateway-proxy-logrotate.timer
[Unit]
Description=Rotate openclaw-gateway-proxy healthcheck log weekly

[Timer]
OnCalendar=weekly
AccuracySec=1h
Persistent=true
Unit=openclaw-gateway-proxy-logrotate.service

[Install]
WantedBy=timers.target
